# Docker/Vedang/supervisor/conf.d/queue.conf
[program:dev_vedanguat_announcement_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_announcement_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_vedanguat_branch_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_branch_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_checkin_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_checkin_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_checkout_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_checkout_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_city_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_city_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_claim_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_claim_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_employee_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_employee_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_geo_data]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_geo_data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_leave_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_leave_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_profile_pic_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_profile_pic_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_regularization_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_regularization_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_shift_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_shift_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_shift_roster]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_shift_roster
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_survey_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_survey_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_user_token]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_user_token
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_weekoff_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_weekoff_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_asset_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_asset_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_shuffle_request_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_shuffle_request_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord.log


[program:dev_vedanguat_shuffle_approval]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_shuffle_approval
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord.log


[program:dev_vedanguat_shuffle_rejected]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_shuffle_rejected
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord.log


[program:dev_vedanguat_performance_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_ikya_performance_upload    --timeout=3600
umask=000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_hierarchy_values_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_hierarchy_values_upload    --timeout=3600
umask=000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_Targets_Import]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_vedanguat_Targets_Import

umask=0000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_ikya_Targets_Import]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5 --queue=dev_ikya_Targets_Import

umask=0000
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_beat_bulk_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_beat_bulk_upload
autostart=true
autorestart=true
redirect_stderr=true
umask=0000

stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_vedanguat_employee_outlet_map]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_employee_outlet_map
autostart=true
autorestart=true
redirect_stderr=true
umask=0000
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log



[program:dev_vedanguat_employee_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5  --queue=dev_vedanguat_employee_notification
autostart=true
autorestart=true
redirect_stderr=true
umask=0000
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_vedanguat_employee_update_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_employee_update_upload
autostart=true
autorestart=true
redirect_stderr=true
umask=0000

stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_vedanguat_survey_one_employee_cache_update]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_survey_one_employee_cache_update
autostart=true
autorestart=true
redirect_stderr=true
umask=0000
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_vedanguat_survey_cache_update]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_survey_cache_update
autostart=true
autorestart=true
redirect_stderr=true
umask=0000
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_posm_inventory_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_posm_inventory_upload
autostart=true
autorestart=true
redirect_stderr=true
umask=0000
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log



[program:dev_vedanguat_send_retail_invoice_mail]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_send_retail_invoice_mail
autostart=true
autorestart=true
redirect_stderr=true
umask=0000
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_medical_certificate_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_medical_certificate_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_farm_certificate_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_farm_certificate_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_fa_certificate_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_fa_certificate_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_medical_certificate_details_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_medical_certificate_details_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_farm_certificate_details_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_farm_certificate_details_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_fa_certificate_details_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_fa_certificate_details_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_pm_report_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_pm_report_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_vedang_asset_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_vedang_asset_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log


[program:dev_vedanguat_exit_mails]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_exit_mails
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_exit_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_exit_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_asset_v2_push_notification]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_asset_v2_push_notification
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_vedang_induction_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5    --queue=dev_vedanguat_vedang_induction_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log

[program:dev_vedanguat_vedang_ohs_audit_upload]
command=php /var/www/html/vhosts/vedang-uat/inedge-backend/artisan queue:work --tries=1 --sleep=5  --queue=dev_vedanguat_vedang_ohs_audit_upload
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/vhosts/vedang-uat/inedge-backend/storage/logs/supervisord.log
